# .config/autostart/warp.sh.desktop
[Desktop Entry]
Exec=/home/nic/Documents/warp.sh
Icon=application-x-shellscript
Name=warp.sh
Type=Application
X-KDE-AutostartScript=true
